# crontab
0 0 */1 * * cd /usr/src/app && /usr/bin/python pod.py fetch && /usr/bin/python pod.py download
